# screenocr-daemon.service
[Unit]
Description=ScreenOCR Hotkey Daemon
Documentation=https://github.com/adikul1023/screenOCR
After=graphical-session-pre.target
PartOf=graphical-session-pre.target

[Service]
Type=simple
ExecStart=%h/.local/bin/screenocr daemon start
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Wayland-specific environment
Environment="DISPLAY="
Environment="WAYLAND_DISPLAY=wayland-0"

# Allow keyboard input
PrivateDevices=no
DevicePolicy=auto
DeviceAllow=/dev/input/event* rw
DeviceAllow=/dev/uinput rw

[Install]
WantedBy=graphical-session-pre.target
